# docs/service/com.remote-coder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.remote-coder</string>
  <key>ProgramArguments</key>
  <array>
    <string>/REPLACE/with/node</string>
    <string>/REPLACE/with/remote-coder/dist/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>REMOTE_CODER_DATA_DIR</key>
    <string>/REPLACE/with/home/.config/remote-coder</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
