# docs/service/remote-coder.service
[Unit]
Description=remote-coder — operate Claude Code sessions remotely
After=network-online.target

[Service]
ExecStart=/REPLACE/with/node /REPLACE/with/remote-coder/dist/index.js
Environment=REMOTE_CODER_DATA_DIR=%h/.config/remote-coder
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
